# cron-systemd/borg-scripts/borg-home.cron
# Cron job for Borg Home Backup
# Runs backup_t14s_home every hour at minute 0
0 * * * * . $HOME/.config/borg/env && $HOME/bin/backup_t14s_home >> $HOME/.cache/borg_home.log 2>&1
